# crontab
0 0 * * * python /app/app.py > /dev/stdout
